# attendance-client-fixed.service
[Unit]
Description=ICカード勤怠管理システム - ラズパイクライアント
After=network.target pcscd.service
Wants=network.target pcscd.service

[Service]
Type=simple
User=raspberry
Group=pcscd
WorkingDirectory=/home/raspberry/Desktop/attendance/card_reader_improved
ExecStartPre=/bin/sleep 5
ExecStartPre=/bin/bash -c 'until systemctl is-active --quiet pcscd; do sleep 1; done'
ExecStartPre=/bin/bash -c 'until [ -S /var/run/pcscd/pcscd.comm ]; do sleep 1; done'
ExecStart=/bin/bash /home/raspberry/Desktop/attendance/card_reader_improved/start_pi.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="LANG=ja_JP.UTF-8"
Environment="LC_ALL=ja_JP.UTF-8"
Environment="PATH=/home/raspberry/Desktop/attendance/card_reader_improved/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/raspberry/Desktop/attendance/card_reader_improved/venv"
Environment="PCSCLITE_CSOCK_NAME=/var/run/pcscd/pcscd.comm"

[Install]
WantedBy=multi-user.target

