# attendance-client.service
[Unit]
Description=ICカード勤怠管理システム - ラズパイクライアント
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/card_reader_improved
ExecStart=/usr/bin/python3 /home/pi/card_reader_improved/pi_client.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 環境変数
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

